# nad.service
[Unit]
Description=REST2NAD Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/rest2nad/nad.py /home/pi/rest2nad/config.yaml

[Install]
WantedBy=multi-user.target
